# dist/digest.supervisor.conf
[program:collector]
command = /root/.nvm/versions/node/v5.1.0/bin/node --use_strict /root/d3-digest/src/collector/collector.js
stderr_logfile = /var/log/supervisor/collector-stderr.log
stdout_logfile = /var/log/supervisor/collector-stdout.log

[program:processor]
command = /root/.nvm/versions/node/v5.1.0/bin/node --use_strict /root/d3-digest/src/processor/processor.js
stderr_logfile = /var/log/supervisor/processor-stderr.log
stdout_logfile = /var/log/supervisor/processor-stdout.log

[program:prefetch]
command = /root/.nvm/versions/node/v5.1.0/bin/node --use_strict /root/d3-digest/src/prefetch/prefetch.js
stderr_logfile = /var/log/supervisor/prefetch-stderr.log
stdout_logfile = /var/log/supervisor/prefetch-stdout.log

[program:web]
command = /root/.nvm/versions/node/v5.1.0/bin/node --use_strict /root/d3-digest/src/web/web.js
stderr_logfile = /var/log/supervisor/web-stderr.log
stdout_logfile = /var/log/supervisor/web-stdout.log
